# evolution-ui.service
[Unit]
Description=Evolution UI - OpenClaw Self-Evolution Engine Web Interface
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/evolution-ui
Environment="PATH=/root/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
